# config/crontab
#This cron will purge expired actions every hour
0 * * * * cd /var/www/html/linkenfest/serverSide && php purgeExpiredActions.php >/dev/null 2>&1
